# projectroot/usr/lib/systemd/system/run-varoverlayfs.mount
[Unit]
Description=Overlay for '/var'
Before=local-fs.target

[Mount]
Where=/run/varoverlayfs
What=tmpfs
Type=tmpfs
Options=size=20%
