# projectroot/usr/lib/systemd/system/tee-supplicant.service
[Unit]
Description=TEE Supplicant

[Service]
Type=simple
ExecStart=/usr/bin/tee-supplicant

[Install]
WantedBy=multi-user.target
